# app-misc/co2mon/files/init.d/co2mond-rrd
#!/sbin/openrc-run
PIDFILE=/var/run/co2mond.pid
start() {
	ebegin "Starting co2mond"
	start-stop-daemon \
		--start \
		--quiet \
		--background \
		--interpreted \
		--exec /usr/bin/co2mond-rrd \
		--pidfile "$PIDFILE" \
		--make-pidfile
	eend $?
}

stop() {
	ebegin "Stopping co2mond"
	kill -- -`cat "$PIDFILE"`
	rm "$PIDFILE"
	eend $?
}
