# net-misc/hamachi/files/tuncfg.initd
#!/sbin/runscript
# Copyright 1999-2006 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

opts="start stop status"

depend() {
        need net
}

start() {
	ebegin "Starting tuncfg"
	start-stop-daemon --start --exec /sbin/tuncfg
	eend $?
}

stop() {
	ebegin "Stopping tuncfg"
	start-stop-daemon --stop --name tuncfg 
	eend $?
}

