# net-misc/wide-dhcpv6/files/dhcp6x.initd.in
#!/sbin/runscript
# Copyright 1999-2010 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

depend() {
	need net
}

start() {
	ebegin "Starting dhcp6x"
	start-stop-daemon --start --quiet --exec /usr/sbin/dhcp6x -- -c /etc/wide-dhcpv6/dhcp6x.conf ${DHCP6X_OPTS}
	eend $?
}

stop() {
	ebegin "Stopping dhcp6x"
	start-stop-daemon --stop --quiet --exec /usr/sbin/dhcp6x
	eend $?
}
