# net-im/pyvk-t/files/initd
#!/sbin/runscript

depend() {
need net
use jabber-server
}

start() {
ebegin "Starting Vkontakte.Ru Jabber Transport"
start-stop-daemon --start --quiet --pidfile /var/run/jabber/pyvk.pid --chuid jabber:jabber \
--exec /usr/bin/python /usr/lib/python2.5/site-packages/pyvk/pyvk-t.py
eend $?
}

stop() {
ebegin "Stopping Vkontakte.Ru Jabber Transport"
start-stop-daemon --stop --quiet --pidfile /var/run/jabber/pyvk.pid \
--exec /usr/bin/killall pyvk-t*
eend $?
}
