# net-im/synapse/files/synapse.initd
#!/sbin/openrc-run
# Copyright 1999-2018 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

user=${SYNAPSE_USER:-synapse}
group=${SYNAPSE_GROUP:-synapse}
pidfile="/var/run/synapse/synapse.pid"
command="/usr/bin/synctl"
command_args="${SYNAPSE_OPTS:-start /etc/synapse}"
start_stop_daemon_args="--user ${user} --group ${group} --pidfile ${pidfile}"
